# systemd/usb-sync.path
[Unit]
Description=Watch for USB mount

[Path]
PathExists=/media/USER/MYBACKUP   # IMPORTANT: CHANGE TO PATH OF YOUR USB DRIVE.

[Install]
WantedBy=default.target
