# systemd/usb-sync.path
[Unit]
Description=Watch for USB mount

[Path]
PathExists=/media/USER/MYBACKUP   # Change to your user

[Install]
WantedBy=default.target
