# systemd/usb-sync.path
[Unit]
Description=Watch for USB mount

[Path]
PathExists=/media/joseph/MYBACKUP

[Install]
WantedBy=default.target
