# systemd/usb-sync.service
[Unit]
Description=USB Auto Sync Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb_sync.py
